# systemd/voice-gradio.service
[Unit]
Description=Voice AI Gradio UI
After=network.target voice-api.service

[Service]
User=voiceai
Group=voiceai
WorkingDirectory=/opt/voice-ai/app
EnvironmentFile=/opt/voice-ai/config/.env
Environment=COQUI_TOS_AGREED=1
ExecStart=/opt/voice-ai/.venv/bin/python -m app.ui.gradio_app
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
